# systemd/adaptive-snowsampler.service
[Unit]
Description=Launch SnowSampler Drone node as a service
After=network.target systemd-udevd.service

[Service]
EnvironmentFile=/home/user/catkin_ws/src/adaptive-snowsampler/systemd/environment.conf
ExecStart=/bin/bash /home/user/catkin_ws/src/adaptive-snowsampler/systemd/run_controller.sh
User=user
Group=user

[Install]
WantedBy=multi-user.target